# config/xinetd/xinetd.d/nc_listen
# default: on
# description: nc_listen
service nc_listen
{
  port		= 3336
  socket_type	= stream
  protocol	= tcp
  wait		= no
  user		= root
  server	= /pxeinstall/httpd/ospcdeploy/nc_listen.sh
  disable	= no
}
